# dot_config/systemd/user/ollama.service
[Unit]
Description=Ollama Service

[Service]
ExecStart=%h/.local/bin/ollama serve
Restart=always

[Install]
WantedBy=default.target
